# com.tourguide.updatecache.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tourguide.updatecache</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/jameslaurie/projects/ride-monitor-backend/updateCache.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/jameslaurie/projects/ride-monitor-backend</string>
    
    <key>StandardOutPath</key>
    <string>/Users/jameslaurie/projects/ride-monitor-backend/cache-update.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jameslaurie/projects/ride-monitor-backend/cache-update-error.log</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
